# .config/systemd/user/wallpaper-change.service
[Unit]
Description=Changes active wallpaper to random one

[Service]
Type=oneshot
ExecStart=/home/arsuhinars/.scripts/wallpapers.sh change

